# crontab
* * * * * cd /var/www/html && /usr/local/bin/php artisan fetch:stock-prices >> /var/log/stock-prices-fetcher-cron.log 2>&1
